# br2external_rigdo/overlay/etc/sv/syslogd/run
#!/bin/sh

exec /sbin/syslogd -n -C2048
