# br2external_rigdo/overlay/etc/inittab
# Startup the system
::sysinit:/etc/init/rc.sysinit

# Start services
::respawn:/usr/bin/runsvdir /var/service
